# com.jesseclaven.appearancenotifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jesseclaven.appearancenotifier</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/AppearanceNotifier</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/appearancenotifier.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/appearancenotifier.out</string>
</dict>
</plist>
